# files/usr/etc/profile.d/40-host-aliases.sh
alias ll='ls -lah'
alias la='ls -A'
alias l='ls -CF'
alias ..='cd ..'
alias ...='cd ../..'
