# files/usr/etc/profile.d/90-welcome-banner.sh
# Show a small welcome banner for interactive shells
case $- in
  *i*)
    echo "Welcome to $(hostname)"
    ;;
esac
